# home/Library/LaunchAgents/com.user.archisteamfarm.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.archisteamfarm</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{ .chezmoi.homeDir }}/Library/Application Support/ArchiSteamFarm/ArchiSteamFarm</string>
    <string>--no-restart</string>
    <string>--service</string>
  </array>
  <key>WorkingDirectory</key>
  <string>{{ .chezmoi.homeDir }}/Library/Application Support/ArchiSteamFarm</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <true/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>StandardOutPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/asf.out.log</string>
  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/asf.err.log</string>
</dict>
</plist>
