# home/private_Library/private_LaunchAgents/private_com.user.obsidian-mirror.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique label for this job -->
    <key>Label</key>
    <string>com.user.obsidian-mirror</string>

    <!-- Path to the sync script — adjust if you put it elsewhere -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>$HOME/.local/bin/sync-obsidian.sh</string>
    </array>

    <!-- Run every 300 seconds (5 minutes) -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Also run immediately when loaded -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/obsidian-mirror.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/obsidian-mirror.log</string>

    <!-- Don't restart on failure — just wait for next interval -->
    <key>KeepAlive</key>
    <false/>

    <!-- Low priority — don't hog resources -->
    <key>Nice</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
